# source/docker-supp/supervisor.cfg
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:webserver]
command=/usr/bin/java -jar /home/ihatestau/webserver.jar -host http://0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:masker]
command=/usr/bin/java -jar /home/ihatestau/masker.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scraper]
command=/usr/bin/java -jar /home/ihatestau/scraper.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:netserver]
command=/home/ihatestau/neural-net-server/serve.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
